# {{cookiecutter.system_project_slug}}/playbook/roles/openresty/templates/profile.d/openresty.sh.j2
export PATH={{openresty_install_dir}}/bin:{{openresty_install_dir}}/nginx/sbin:$PATH